# api/supervisor.conf
[supervisord]
nodaemon=true

[program:app1]
directory=/home/api
command=npm run dev
environment=PORT="8080"
stdout_logfile=/logs/app1.log
stdout_logfile_maxbytes=1GB
redirect_stderr=true